# Installs/25b82B1/Lx/bin/ezconnect.service
[Unit]
Description=EZConnect service
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=1
User=eztechlx
ExecStart=/usr/ezlx/bin/ezconnect -s
LimitRTPRIO=infinity


[Install]
WantedBy=multi-user.target
